# albumweb.service
[Unit]
Description=Album Player Web
After=network.target

[Service]
WorkingDirectory=/home/dyonak/Album-Player-Standalone
ExecStart=/home/dyonak/Album-Player-Standalone/myvenv/bin/python3 Webapp.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
